# deploy/notamax.service
[Unit]
Description=notamax (Node + Socket.io + static)
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/notamax/server
Environment=NODE_ENV=production
Environment=PORT=3001
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
