# deploy/notamax.service
[Unit]
Description=notamax (Node + Socket.io + static)
After=network.target

[Service]
Type=simple
User=notamax
Group=notamax
WorkingDirectory=/opt/notamax/server
Environment=NODE_ENV=production
Environment=PORT=3010
# Если фронт и API с одного хоста — ALLOWED_ORIGINS можно не задавать.
# Для жёсткого списка разрешённых origin:
# Environment=ALLOWED_ORIGINS=https://your-domain.example
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
